# trunk/deskutils/beagle/pkg-plist
bin/beagle-config
bin/beagle-exercise-file-system
bin/beagle-extract-content
bin/beagle-imlogviewer
bin/beagle-index-info
bin/beagle-index-url
bin/beagle-info
bin/beagle-ping
bin/beagle-query
bin/beagle-search
bin/beagle-settings
bin/beagle-shutdown
bin/beagle-status
bin/beagled
etc/beagle/crawl-applications
etc/beagle/crawl-documentation
etc/beagle/crawl-windows
etc/beagle/external-filters.xml.sample
etc/cron.daily/beagle-crawl-system
etc/xdg/autostart/beagle-search-autostart.desktop
etc/xdg/autostart/beagled-autostart.desktop
include/libbeagle/beagle/beagle-client.h
include/libbeagle/beagle/beagle-daemon-information-request.h
include/libbeagle/beagle/beagle-daemon-information-response.h
include/libbeagle/beagle/beagle-empty-response.h
include/libbeagle/beagle/beagle-error-response.h
include/libbeagle/beagle/beagle-finished-response.h
include/libbeagle/beagle/beagle-hit.h
include/libbeagle/beagle/beagle-hits-added-response.h
include/libbeagle/beagle/beagle-hits-subtracted-response.h
include/libbeagle/beagle/beagle-indexable.h
include/libbeagle/beagle/beagle-indexing-service-request.h
include/libbeagle/beagle/beagle-indexing-status-response.h
include/libbeagle/beagle/beagle-informational-messages-request.h
include/libbeagle/beagle/beagle-property.h
include/libbeagle/beagle/beagle-query-part-date.h
include/libbeagle/beagle/beagle-query-part-human.h
include/libbeagle/beagle/beagle-query-part-or.h
include/libbeagle/beagle/beagle-query-part-property.h
include/libbeagle/beagle/beagle-query-part-text.h
include/libbeagle/beagle/beagle-query-part-wildcard.h
include/libbeagle/beagle/beagle-query-part.h
include/libbeagle/beagle/beagle-query.h
include/libbeagle/beagle/beagle-queryable-status.h
include/libbeagle/beagle/beagle-request.h
include/libbeagle/beagle/beagle-response.h
include/libbeagle/beagle/beagle-scheduler-information.h
include/libbeagle/beagle/beagle-search-term-response.h
include/libbeagle/beagle/beagle-shutdown-request.h
include/libbeagle/beagle/beagle-snippet-request.h
include/libbeagle/beagle/beagle-snippet-response.h
include/libbeagle/beagle/beagle-timestamp.h
include/libbeagle/beagle/beagle-util.h
include/libbeagle/beagle/beagle.h
lib/beagle/Backends/EvolutionBackends.dll
lib/beagle/Backends/EvolutionBackends.dll.mdb
lib/beagle/Beagle.dll
lib/beagle/Beagle.dll.mdb
lib/beagle/BeagleDaemon.exe
lib/beagle/BeagleDaemon.exe.mdb
lib/beagle/BeagleDaemonLib.dll
lib/beagle/BeagleDaemonLib.dll.mdb
lib/beagle/BeagleDaemonPlugins.dll
lib/beagle/BeagleDaemonPlugins.dll.mdb
lib/beagle/BuildIndex.exe
lib/beagle/BuildIndex.exe.mdb
lib/beagle/Config.exe
lib/beagle/Config.exe.mdb
lib/beagle/DumpIndex.exe
lib/beagle/DumpIndex.exe.mdb
lib/beagle/ExerciseFileSystem.exe
lib/beagle/ExerciseFileSystem.exe.mdb
lib/beagle/ExtractContent.exe
lib/beagle/ExtractContent.exe.mdb
lib/beagle/Filters/Filters.dll
lib/beagle/Filters/Filters.dll.mdb
lib/beagle/ImLogViewer.exe
lib/beagle/ImLogViewer.exe.mdb
lib/beagle/Images.dll
lib/beagle/Images.dll.mdb
lib/beagle/IndexHelper.exe
lib/beagle/IndexHelper.exe.mdb
lib/beagle/IndexWebContent.exe
lib/beagle/IndexWebContent.exe.mdb
lib/beagle/Info.exe
lib/beagle/Info.exe.mdb
lib/beagle/ManageIndex.exe
lib/beagle/ManageIndex.exe.mdb
lib/beagle/Query.exe
lib/beagle/Query.exe.mdb
lib/beagle/Search.exe
lib/beagle/Search.exe.mdb
lib/beagle/Settings.exe
lib/beagle/Settings.exe.mdb
lib/beagle/Shutdown.exe
lib/beagle/Shutdown.exe.mdb
lib/beagle/UiUtil.dll
lib/beagle/UiUtil.dll.config
lib/beagle/UiUtil.dll.mdb
lib/beagle/Util.dll
lib/beagle/Util.dll.config
lib/beagle/Util.dll.mdb
lib/beagle/beagled-index-helper
lib/beagle/libbeagleglue.a
lib/beagle/libbeagleglue.la
lib/beagle/libbeagleglue.so
lib/beagle/libbeagleglue.so.0
lib/beagle/libbeagleuiglue.a
lib/beagle/libbeagleuiglue.la
lib/beagle/libbeagleuiglue.so
lib/beagle/libbeagleuiglue.so.0
%%WV:%%lib/beagle/libwv1glue.a
%%WV:%%lib/beagle/libwv1glue.so
%%WV:%%lib/beagle/libwv1glue.so.0
%%EPIPHANY:%%lib/epiphany/2.16/extensions/beagle.xml
%%EPIPHANY:%%lib/epiphany/2.16/extensions/libbeagleextension.a
%%EPIPHANY:%%lib/epiphany/2.16/extensions/libbeagleextension.la
%%EPIPHANY:%%lib/epiphany/2.16/extensions/libbeagleextension.so
lib/libbeagle.a
lib/libbeagle.la
lib/libbeagle.so
lib/libbeagle.so.0
libdata/pkgconfig/beagle-0.0.pc
libdata/pkgconfig/beagle-daemon.pc
libdata/pkgconfig/beagle-ui-0.0.pc
libdata/pkgconfig/libbeagle-0.0.pc
sbin/beagle-build-index
sbin/beagle-dump-index
sbin/beagle-manage-index
share/applications/beagle-search.desktop
share/applications/beagle-settings.desktop
%%DATADIR%%/beagle.xpi
%%DOCSDIR%%/BeagleClient.html
%%DOCSDIR%%/BeagleDaemonInformationResponse.html
%%DOCSDIR%%/BeagleEmptyResponse.html
%%DOCSDIR%%/BeagleErrorResponse.html
%%DOCSDIR%%/BeagleFinishedResponse.html
%%DOCSDIR%%/BeagleHitsAddedResponse.html
%%DOCSDIR%%/BeagleHitsSubtractedResponse.html
%%DOCSDIR%%/BeagleIndexingServiceRequest.html
%%DOCSDIR%%/BeagleQuery.html
%%DOCSDIR%%/BeagleRequest.html
%%DOCSDIR%%/BeagleResponse.html
%%DOCSDIR%%/BeagleShutdownRequest.html
%%DOCSDIR%%/BeagleSnippetRequest.html
%%DOCSDIR%%/BeagleSnippetResponse.html
%%DOCSDIR%%/beagle-BeagleCancelledResponse.html
%%DOCSDIR%%/beagle-BeagleDaemonInformationRequest.html
%%DOCSDIR%%/beagle-BeagleHit.html
%%DOCSDIR%%/beagle-beagle-indexable.html
%%DOCSDIR%%/beagle-beagle-property.html
%%DOCSDIR%%/beagle-beagle-query-part.html
%%DOCSDIR%%/beagle-beagle-timestamp.html
%%DOCSDIR%%/beagle.devhelp
%%DOCSDIR%%/beagle.devhelp2
%%DOCSDIR%%/ch01.html
%%DOCSDIR%%/ch02.html
%%DOCSDIR%%/ch03.html
%%DOCSDIR%%/ch04.html
%%DOCSDIR%%/ch05.html
%%DOCSDIR%%/home.png
%%DOCSDIR%%/index.html
%%DOCSDIR%%/index.sgml
%%DOCSDIR%%/left.png
%%DOCSDIR%%/right.png
%%DOCSDIR%%/style.css
%%DOCSDIR%%/up.png
share/locale/ar/LC_MESSAGES/beagle.mo
share/locale/bg/LC_MESSAGES/beagle.mo
share/locale/ca/LC_MESSAGES/beagle.mo
share/locale/cs/LC_MESSAGES/beagle.mo
share/locale/da/LC_MESSAGES/beagle.mo
share/locale/de/LC_MESSAGES/beagle.mo
share/locale/dz/LC_MESSAGES/beagle.mo
share/locale/el/LC_MESSAGES/beagle.mo
share/locale/en_CA/LC_MESSAGES/beagle.mo
share/locale/en_GB/LC_MESSAGES/beagle.mo
share/locale/es/LC_MESSAGES/beagle.mo
share/locale/fi/LC_MESSAGES/beagle.mo
share/locale/fr/LC_MESSAGES/beagle.mo
share/locale/gl/LC_MESSAGES/beagle.mo
share/locale/he/LC_MESSAGES/beagle.mo
share/locale/hi/LC_MESSAGES/beagle.mo
share/locale/hu/LC_MESSAGES/beagle.mo
share/locale/it/LC_MESSAGES/beagle.mo
share/locale/ja/LC_MESSAGES/beagle.mo
share/locale/ka/LC_MESSAGES/beagle.mo
share/locale/ko/LC_MESSAGES/beagle.mo
share/locale/lt/LC_MESSAGES/beagle.mo
share/locale/lv/LC_MESSAGES/beagle.mo
share/locale/mk/LC_MESSAGES/beagle.mo
share/locale/nb/LC_MESSAGES/beagle.mo
share/locale/nl/LC_MESSAGES/beagle.mo
share/locale/pa/LC_MESSAGES/beagle.mo
share/locale/pl/LC_MESSAGES/beagle.mo
share/locale/pt_BR/LC_MESSAGES/beagle.mo
share/locale/ru/LC_MESSAGES/beagle.mo
share/locale/sl/LC_MESSAGES/beagle.mo
share/locale/sr/LC_MESSAGES/beagle.mo
share/locale/sr@Latn/LC_MESSAGES/beagle.mo
share/locale/sv/LC_MESSAGES/beagle.mo
share/locale/th/LC_MESSAGES/beagle.mo
share/locale/tr/LC_MESSAGES/beagle.mo
share/locale/uk/LC_MESSAGES/beagle.mo
share/locale/vi/LC_MESSAGES/beagle.mo
share/locale/zh_CN/LC_MESSAGES/beagle.mo
share/locale/zh_HK/LC_MESSAGES/beagle.mo
share/locale/zh_TW/LC_MESSAGES/beagle.mo
share/pixmaps/system-search.png
@dirrmtry share/locale/zh_TW/LC_MESSAGES
@dirrmtry share/locale/zh_TW
@dirrmtry share/locale/zh_HK/LC_MESSAGES
@dirrmtry share/locale/zh_HK
@dirrmtry share/locale/zh_CN/LC_MESSAGES
@dirrmtry share/locale/zh_CN
@dirrmtry share/locale/uk/LC_MESSAGES
@dirrmtry share/locale/uk
@dirrmtry share/locale/sr/LC_MESSAGES
@dirrmtry share/locale/sr
@dirrmtry share/locale/sl/LC_MESSAGES
@dirrmtry share/locale/sl
@dirrmtry share/locale/pl/LC_MESSAGES
@dirrmtry share/locale/pl
@dirrmtry share/locale/ko/LC_MESSAGES
@dirrmtry share/locale/ko
@dirrmtry share/locale/gl/LC_MESSAGES
@dirrmtry share/locale/gl
@dirrmtry share/locale/fi/LC_MESSAGES
@dirrmtry share/locale/fi
@dirrmtry share/locale/el/LC_MESSAGES
@dirrmtry share/locale/el
@dirrmtry share/locale/dz/LC_MESSAGES
@dirrmtry share/locale/dz
@dirrmtry share/locale/cs/LC_MESSAGES
@dirrmtry share/locale/cs
@dirrmtry share/locale/ca/LC_MESSAGES
@dirrmtry share/locale/ca
@dirrm %%DOCSDIR%%
@dirrm %%DATADIR%%
%%EPIPHANY:%%@dirrmtry lib/epiphany/2.16/extensions
%%EPIPHANY:%%@dirrmtry lib/epiphany/2.16
%%EPIPHANY:%%@dirrmtry lib/epiphany
@dirrm lib/beagle/Filters
@dirrm lib/beagle/Backends
@dirrm lib/beagle
@dirrm include/libbeagle/beagle
@dirrm include/libbeagle
@dirrmtry etc/xdg/autostart
@dirrmtry etc/xdg
@dirrmtry etc/cron.daily
@dirrm etc/beagle
